# systemd/dashboard.service
[Unit]
Description=E-ink dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-/etc/default/eink-dashboard
Environment=PYTHONUNBUFFERED=1
ExecStart=__INSTALL_DIR__/.venv/bin/python -m eink_dashboard.cli
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
